# amdgpu-fan
#!/sbin/openrc-run

name="Fan controller for amdgpus"
command="/usr/bin/amdgpu-fan"
pidfile="/var/run/amdgpu-fan.pid"
command_background=true

